# site/profiles/files/puppetmodule-daily.timer
# THIS FILE IS MANAGED BY PUPPET
[Unit]
Description=partial module update each day

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target
